# setup/dbt-docs.service
[Unit]
Description=dbt Documentation Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=amit
Group=datausers
WorkingDirectory=/srv/projects/dbt-asc
Environment="PATH=/home/amit/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/amit/.snowflake_env
ExecStartPre=/home/amit/.local/bin/dbt docs generate
ExecStart=/home/amit/.local/bin/dbt docs serve --port 8082
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dbt-docs

[Install]
WantedBy=multi-user.target
